# infra/systemd/a9-recovery-loop.service
[Unit]
Description=A9 recovery planning loop
After=network-online.target a9-control-api.service

[Service]
Type=simple
WorkingDirectory=/root/a9
EnvironmentFile=-/root/a9/.env
ExecStart=/root/a9/scripts/a9_recovery_loop.py --controller-url http://127.0.0.1:8787 --interval-seconds 60 --timeout 10 --max-actions 3
Restart=always
RestartSec=10
KillSignal=SIGINT

[Install]
WantedBy=multi-user.target
